# core/startup/affilly.service
[Unit]
Description=AutoAffiliateHub-X2 Service
After=network.target
Wants=network.target

[Service]
Type=simple
User=affilly
Group=affilly
WorkingDirectory=/opt/affilly/AutoAffiliateHub-X2
Environment=NODE_ENV=production
Environment=PYTHONPATH=/opt/affilly/AutoAffiliateHub-X2
Environment=AFFILLY_LOG_LEVEL=INFO
ExecStart=/usr/bin/python3 /opt/affilly/AutoAffiliateHub-X2/core/startup/orchestrator.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=10
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Resource limits
MemoryHigh=1G
MemoryMax=2G
CPUWeight=100
TasksMax=100

# Security
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/affilly/AutoAffiliateHub-X2/logs
ReadWritePaths=/opt/affilly/AutoAffiliateHub-X2/core/startup/logs
ReadWritePaths=/tmp

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=affilly

[Install]
WantedBy=multi-user.target